# build/supervisor/conf.d/zabbix.conf
[program:zabbix]
command = /usr/sbin/zabbix_agentd -f
user = root
autostart = true
autorestart = false